# files/update-motd.d/20_Welcome_Message
#! /bin/sh
logger "Setting MOTD with $0"
# Following text thanks to http://patorjk.com/software/taag/
echo '======================================================================================'
echo '          _______  _        _______  _______  _______  _______      _________ _______ '
echo '|\     /|(  ____ \( \      (  ____ \(  ___  )(       )(  ____ \     \__   __/(  ___  )'
echo '| )   ( || (    \/| (      | (    \/| (   ) || () () || (    \/        ) (   | (   ) |'
echo '| | _ | || (__    | |      | |      | |   | || || || || (__            | |   | |   | |'
echo '| |( )| ||  __)   | |      | |      | |   | || |(_)| ||  __)           | |   | |   | |'
echo '| || || || (      | |      | |      | |   | || |   | || (              | |   | |   | |'
echo '| () () || (____/\| (____/\| (____/\| (___) || )   ( || (____/\        | |   | (___) |'
echo '(_______)(_______/(_______/(_______/(_______)|/     \|(_______/        )_(   (_______)'
echo '                                                                                      '
echo ' _                 _______         _______  _______  _______                  _       '
echo '( \      |\     /|(  ____ \       (  ___  )(  ____ )(  ____ \       |\     /|| \    /\'
echo '| (      | )   ( || (    \/       | (   ) || (    )|| (    \/       | )   ( ||  \  / /'
echo '| |      | |   | || |             | |   | || (____)|| |             | |   | ||  (_/ / '
echo '| |      | |   | || | ____        | |   | ||     __)| | ____        | |   | ||   _ (  '
echo '| |      | |   | || | \_  )       | |   | || (\ (   | | \_  )       | |   | ||  ( \ \ '
echo '| (____/\| (___) || (___) |   _   | (___) || ) \ \__| (___) |   _   | (___) ||  /  \ \'
echo '(_______/(_______)(_______)  (_)  (_______)|/   \__/(_______)  (_)  (_______)|_/    \/'
echo '======================================================================================'
echo '                                                                                      '
