# files/update-motd.d/25_Inject_motd_tail
#! /bin/sh
logger "Setting MOTD with $0"
[ ! -f /etc/motd.tail ] && exit 0
cat /etc/motd.tail
echo ''
