# files/update-motd.d/30_Admin_Contact
#! /bin/sh
logger "Setting MOTD with $0"
echo "    If you have any issues with this host, please contact mailto:admin@lug.org.uk.    "
echo ''
